# overlays/module-v4l2loopback/built/usr/lib/systemd/system/dslr-webcam.service
[Unit]
Description=DSLR webcam script using gphoto2 and ffmpeg

[Service]
ExecStart=/usr/local/bin/dslr-webcam --first
Restart=always
RestartSec=10

[Install]
WantedBy=basic.target
