# overlays/game-path-of-exile-2/usr/lib/systemd/user/neversink-update.service
[Unit]
PartOf=graphical-session.target
After=graphical-session.target mako.service
Requisite=graphical-session.target
Requires=mako.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/neversink-update.sh
